# crontab
0 */6 * * * /root/backup.sh > /var/log/db_backup.log
